# scripts/app-atelier.service.template
[Unit]
Description=App Atelier Informatique (Next.js)
After=network.target

[Service]
Type=simple
User={{APP_USER}}
WorkingDirectory={{APP_DIR}}
Environment=NODE_ENV=production
Environment=PORT={{PORT}}
EnvironmentFile={{APP_DIR}}/.env
ExecStart={{NPM_BIN}} run start -- -p {{PORT}}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
